# config/omarchy-mirror-upload-stable.timer
[Unit]
Description=Omarchy Mirror Upload Stable — Last Monday @ 10am UTC of Each Month

[Timer]
OnCalendar=Mon *-*-24..31 10:00
RandomizedDelaySec=1800
Persistent=true

[Install]
WantedBy=timers.target
